# packaging/debian_amd64/lib/systemd/system/data-warehouse.path
[Unit]
Description=openbank data warehouse namespace
PartOf=data-warehouse.service
ConditionPathExists=/etc/init/data-warehouse.conf

[Path]
PathChanged=/etc/init/data-warehouse.conf
PathModified=/etc/init/data-warehouse.conf
Unit=data-warehouse.service

[Install]
WantedBy=data-warehouse.service
