# packaging/debian/data-warehouse-watcher.path
[Unit]
Description=openbank data warehouse configuration watcher
After=data-warehouse.service

[Path]
PathChanged=/etc/data-warehouse/conf.d
Unit=data-warehouse-watcher.service

[Install]
WantedBy=paths.target
