# packaging/debian/data-warehouse.path
[Unit]
Description=openbank data warehouse namespace
PartOf=data-warehouse.service
ConditionPathExists=/etc/data-warehouse/conf.d/init.conf

[Path]
PathChanged=/etc/data-warehouse/conf.d/init.conf
PathModified=/etc/data-warehouse/conf.d/init.conf
Unit=data-warehouse.service

[Install]
WantedBy=data-warehouse.service
