# packaging/debian/data-warehouse-watcher.path
[Unit]
Description=openbank data warehouse configuration watcher

[Path]
PathChanged=/etc/data-warehouse/conf.d/init.conf
Unit=data-warehouse-watcher.service

[Install]
WantedBy=multi-user.target